# deploy/discord-bot-watchdog.service
[Unit]
Description=Check discord-bot.service is running, restart it if not
After=discord-bot.service

[Service]
# Runs as root (no User= override): reset-failed/start on another system
# unit needs privileges a plain user account doesn't have.
Type=oneshot
ExecStart=/home/pi/discord-music-bot/deploy/watchdog_discord_bot.sh
