# deploy/yt-dlp-update.timer
[Unit]
Description=Weekly yt-dlp update check

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=1h

[Install]
WantedBy=timers.target
